# Ansible-Project/mops/roles/web/templates/flask.conf
description "flask"
start on stopped rc RUNLEVEL=[2345]
respawn
exec python3 /home/app.py
